# config/com.usblocker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.usblocker.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/python3</string>
        <string>/opt/usb-blocker/service/macos_daemon.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/usb-blocker</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/usb-blocker/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/usb-blocker/launchd.err.log</string>

    <key>UserName</key>
    <string>root</string>
</dict>
</plist>
